# supervisor/kafka.ini
[program:kafka]
command=/usr/bin/start-kafka.sh
stdout_logfile=/var/log/kafka.out
stderr_logfile=/var/log/kafka.err
autostart=true
autorestart=true