# supervisor/zookeeper.ini
[program:zookeeper]
command=zkServer.sh start-foreground
stdout_logfile=/var/log/zookeeper.out
stderr_logfile=/var/log/zookeeper.err
autostart=true
autorestart=true